# gui-libs/greetd/files/greetd
#!/sbin/openrc-run

name="greetd"
command="/usr/bin/greetd"
command_background="yes"
description="greetd login daemon"
